# root-stow/vmware-shared-dirs/etc/systemd/system/vmshare.automount
[Unit]
Description=Load VMware shared folders
ConditionPathExists=.host:/
ConditionVirtualization=vmware

[Automount]
Where=/vmshare/

[Install]
WantedBy=multi-user.target
